# examples/systemd/keyvault-certdeploy-periodic.timer.j2
[Unit]
Description=Periodically check for certificate updates

[Timer]
OnBootSec=1d
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
